# deploy/supervisor/exchange-server/mbta-exchange.conf
[program:mbta-exchange]
command=/opt/mbta-agentcy/venv/bin/python -m uvicorn exchange_agent.exchange_server:app --host 0.0.0.0 --port 8100
directory=/opt/mbta-agentcy/src
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/mbta-exchange.out.log
stderr_logfile=/var/log/mbta-exchange.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=
    PYTHONPATH="/opt/mbta-agentcy/src",
    REGISTRY_URL="http://97.107.132.213:6900",
    ANS_ENABLED="true",
    ANS_RESOLVER_URL="http://50.116.53.133:8200",
    ANS_RESOLVER_TIMEOUT="5.0",
    ANS_TLD="agents.dataworksai.com",
    ANS_APP="mbta-transit-ci",
    USE_SLIM="true",
    SLIM_ENDPOINT="http://96.126.111.107:46357",
    SLIM_SHARED_SECRET="%(ENV_SLIM_SHARED_SECRET)s",
    SLIM_ORG="mbta",
    SLIM_NS="transit-ci",
    SLIM_TIMEOUT="60.0",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    MBTA_API_KEY="%(ENV_MBTA_API_KEY)s"
