# deploy/supervisor/agents-server/mbta-slim-stopfinder.conf
[program:mbta-slim-stopfinder]
; StopFinder agent SLIM transport (runs alongside mbta-stopfinder on :50053)
command=/opt/mbta-agents/venv/bin/python /opt/mbta-agents/agents/stopfinder/run_slim.py
directory=/opt/mbta-agents/agents/stopfinder
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/mbta-slim-stopfinder.out.log
stderr_logfile=/var/log/mbta-slim-stopfinder.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=
    PYTHONPATH="/opt/mbta-agents",
    AGENT_HOST="96.126.111.107",
    SLIM_LABEL="stopfinder",
    SLIM_ENDPOINT="http://96.126.111.107:46357",
    SLIM_SHARED_SECRET="%(ENV_SLIM_SHARED_SECRET)s",
    SLIM_ORG="mbta",
    SLIM_NS="transit-ci",
    MBTA_API_KEY="%(ENV_MBTA_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s"
