# deploy/supervisor/agents-server/mbta-slim-controller.conf
[program:mbta-slim-controller]
; Runs the Cisco AGNTCY SLIM controller via Docker.
; All agents + exchange connect to this on port 46357.
command=docker run --rm --name slim-controller
    -v /opt/mbta-agents/slim/config.yaml:/config.yaml
    -p 46357:46357
    ghcr.io/agntcy/slim:1.0.0
    /slim --config /config.yaml
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/mbta-slim-controller.out.log
stderr_logfile=/var/log/mbta-slim-controller.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; No environment block — docker handles its own env
